# deploy/systemd/vectra-vps-disk-cleanup.service
[Unit]
Description=Vectra conservative VPS disk cleanup
Wants=docker.service
After=docker.service
ConditionPathExists=/opt/vectra-prorouter/deploy/scripts/vps-disk-cleanup.sh

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash /opt/vectra-prorouter/deploy/scripts/vps-disk-cleanup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=yes
